# ansible/playbooks/templates/bytefreezer-proxy.service.j2
[Unit]
Description=ByteFreezer Proxy - UDP Data Collection and Forwarding Service
Documentation=https://github.com/{{ github_repo }}
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ service.user }}
Group={{ service.group }}
WorkingDirectory={{ paths.binary_dir }}
ExecStart={{ paths.binary_path }} --config="{{ paths.config_path }}"
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Restart settings
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging to files
StandardOutput=append:{{ paths.log_dir }}/{{ service.name }}.log
StandardError=append:{{ paths.log_dir }}/{{ service.name }}-error.log

# Environment
Environment=CONFIG_FILE={{ paths.config_path }}
Environment=LOG_LEVEL={{ config.logging.level | upper }}

[Install]
WantedBy=multi-user.target